# scripts/backup/nunux/cron.d/backup-nunux
# Backup Nunux archives to remote archive server
00 06 * * * root [ -x /usr/local/sbin/backup-nunux ] && /usr/local/sbin/backup-nunux
