# scripts/backup/home/cron.d/backup-private
# Backup Photos and Videos to remote archive server
00 09 * * * root [ -x /usr/local/sbin/backup-private ] && /usr/local/sbin/backup-private
