# scripts/backup/home/cron.d/backup-nunux
# Backup Nunux archives on this computer
00 10 * * * root [ -x /usr/local/sbin/backup-nunux ] && /usr/local/sbin/backup-nunux
