# scripts/backup/home/cron.d/backup-home
# Backup Home archives to remote archive server
00 07 * * * root [ -x /usr/local/sbin/backup-home ] && /usr/local/sbin/backup-home
